# pkg/lockstep.service.in
[Unit]
Description=Record process resource consumption

[Service]
Nice=19
EnvironmentFile=-/@sysconfdir@/sysconfig/@name@
ExecStart=@prefix@/@bindir@/@name@ $LOCKSTEP_ARGS
StandardOutput=file:@localstatedir@/log/@name@/@name@.log
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
